# source/usr/local/emhttp/plugins/dwpar2cron/defaults/default.cron
# dwpar2cron:
* * * * * /usr/bin/par2cron-cron > /dev/null 2>&1
